# contrib/freebsd/tubenews
#!/bin/sh
#
# PROVIDE: tubenews
# REQUIRE: NETWORKING
# KEYWORD: shutdown
#
# FreeBSD rc.d service for the TubeNews web server.
#
# Install:
#   cp contrib/freebsd/tubenews /usr/local/etc/rc.d/tubenews
#   chmod +x /usr/local/etc/rc.d/tubenews
#
# Then add to /etc/rc.conf (adjust paths):
#   tubenews_enable="YES"
#   tubenews_dir="/usr/local/tubenews"
#
# Optional rc.conf settings and their defaults:
#   tubenews_user="www"                        # system user to run as
#   tubenews_https="NO"                        # YES when behind an HTTPS proxy
#   tubenews_logfile="/var/log/tubenews.log"   # stdout + stderr destination
#
# Manage:
#   service tubenews start | stop | restart | status

. /etc/rc.subr

name="tubenews"
rcvar="${name}_enable"

load_rc_config "$name"

: ${tubenews_enable:="NO"}
: ${tubenews_dir:=""}
: ${tubenews_https:="NO"}
: ${tubenews_user:="www"}
: ${tubenews_logfile:="/var/log/tubenews.log"}

pidfile="/var/run/${name}.pid"

start_cmd="${name}_start"
stop_cmd="${name}_stop"
status_cmd="${name}_status"
restart_cmd="${name}_restart"
extra_commands="status"

_require_dir() {
    if [ -z "${tubenews_dir}" ]; then
        err 1 "Set tubenews_dir in /etc/rc.conf to the TubeNews installation directory."
    fi
    if [ ! -f "${tubenews_dir}/serve.sh" ]; then
        err 1 "serve.sh not found in ${tubenews_dir}"
    fi
}

tubenews_start() {
    _require_dir
    local env_extra=""
    [ "${tubenews_https}" = "YES" ] && env_extra="TUBENEWS_HTTPS=true"
    echo "Starting ${name}."
    /usr/sbin/daemon \
        -p "${pidfile}" \
        -u "${tubenews_user}" \
        -o "${tubenews_logfile}" \
        /usr/bin/env ${env_extra} sh "${tubenews_dir}/serve.sh"
    # To auto-restart on crash, replace the above with:
    #   /usr/sbin/daemon -r -R 5 -p "${pidfile}" -u "${tubenews_user}" \
    #       -o "${tubenews_logfile}" \
    #       /usr/bin/env ${env_extra} sh "${tubenews_dir}/serve.sh"
}

tubenews_stop() {
    if [ -f "${pidfile}" ]; then
        local pid
        pid=$(cat "${pidfile}")
        echo "Stopping ${name} (pid ${pid})."
        kill "${pid}" 2>/dev/null
        wait_for_pids "${pid}"
        rm -f "${pidfile}"
    else
        echo "${name} is not running."
    fi
}

tubenews_status() {
    if [ -f "${pidfile}" ] && kill -0 "$(cat "${pidfile}")" 2>/dev/null; then
        echo "${name} is running as pid $(cat "${pidfile}")."
    else
        echo "${name} is not running."
    fi
}

tubenews_restart() {
    tubenews_stop
    tubenews_start
}

run_rc_command "$1"
